# deploy/sub2api-distributor.service
[Unit]
Description=Sub2API Distributor Deployment Stack
Documentation=https://github.com/your-org/sub2api-distributor
After=network.target docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/sub2api-distributor/deploy
ExecStart=/usr/bin/docker compose up -d --build
ExecStop=/usr/bin/docker compose down
ExecReload=/usr/bin/docker compose up -d --build
TimeoutStartSec=0
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
